# RUN/eventlogger.conf
[program:eventlogger]
command=/root/temp/eventlogger/python/bin/python /root/temp/eventlogger/eventlogger/project/manage.py runserver 0.0.0.0:8080 --insecure
autostart=true
autorestart=true
